# roles/pterodactyl_panel/files/pterosched.timer
[Unit]
Description=Run Artisan scheduled jobs for Pterodactyl Panel

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
